# linux/systemd/cache-cleanup.timer
[Unit]
Description=Run dev cache cleanup weekly

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
