# linux/systemd/zsh-history-cleanup.timer
[Unit]
Description=Run zsh history cleanup weekly

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
